# extensions/gh-runner/koken-runner.service
[Unit]
Description=koken-akademi ephemeral GitHub Actions runner (Docker-isolated, owner-only)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=kokenrunner
# Wrapper root'a-ait sabit yoldan koşar (setup scripti kurar). /opt/linux-ai-server
# kurulu deployment'ta app-kullanıcısına chown'lu olduğundan oradan koşarsak app-user
# wrapper'ı değiştirip kokenrunner (docker-grup ~ root) olarak kod çalıştırabilirdi (Codex :10).
ExecStart=/usr/local/lib/koken-runner/run-ephemeral-loop.sh
Restart=always
RestartSec=15
# StartLimitBurst YOK: normal job-cycle'lar loop-içi iterasyon (servis-restart DEĞİL) →
# eski guard (StartLimitBurst=5) normal yükü kırıyordu; bu tasarımda restart yalnız
# orkestratör-crash'inde olur (#1175 Codex :133 fix).
# KillMode varsayılan (control-group): durdurma/restart'ta child-process'ler de temizlenir
# (eski KillMode=process child sızdırıyordu, Codex :143). Container zaten --rm.

[Install]
WantedBy=multi-user.target
