# automation/crontab
# Linux-AI Server Automation
*/10 * * * * /opt/linux-ai-server/automation/health-check.sh
*/5 * * * * /opt/linux-ai-server/automation/alert-check.sh
0 3 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh daily-backup /opt/linux-ai-server/automation/daily-backup.sh
# Daily test runner — tüm projeler
0 6 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh test-runner /opt/linux-ai-server/automation/test-runner.sh
# Haftalık dependency audit — Pazartesi 09:00
0 9 * * 1 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh weekly-audit /opt/linux-ai-server/automation/weekly-audit.sh
# Panola ERP — Gece demo reset + E2E test
0 4 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh demo-reset-test /opt/linux-ai-server/automation/demo-reset-test.sh
# Sosyal Medya Otomasyon
# E2E Live Site Tests — günlük 07:00
0 7 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh e2e-live-test /opt/linux-ai-server/automation/e2e-live-test.sh
# Günlük operasyon dijesti -> Telegram (08:00, sabah işleri sonrası; NOTHING_NEW guard) [F-C1]
0 8 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh digest /opt/linux-ai-server/automation/digest-send.sh
# Cross-project PR-review poller — her 2 saat, ENABLED 2026-06-04 (spawn+post ACIK; pilot=claude-server; cap 5-run+10-gun; Codex-once-force). Rollback: env-prefix kaldir.
0 */2 * * * DRY_RUN=0 PR_REVIEW_ENABLED=1 SPAWN_ENABLED=1 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh pr-review-poll /opt/linux-ai-server/automation/pr-review-poller.sh
# Klipper HDD backup — VPS Dokploy + volumes (7 gun retention)
0 4 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh pull-vps-backup /opt/linux-ai-server/automation/pull-vps-backup.sh
# One-shot: gozlem haftasi sonu hatirlatma (2026-05-21, self-cleaning)
0 9 21 5 * /opt/linux-ai-server/scripts/observation-week-end-reminder.sh > /var/log/linux-ai-server/observation-reminder.log 2>&1
# Bilge Arena failover — Dokploy down ise Vercel e gecis (3dk kontrol)
# Hafiza maintenance: 90+ gun eski okunmamis discoveries arsivle
30 2 * * * bash /opt/linux-ai-server/scripts/memory-archive-stale.sh
# Hafiza triage: 60+ gun active learning/workaround obsolete (gunluk 03:15)
15 3 * * * bash /opt/linux-ai-server/scripts/memory-triage.sh
# LLM-tabanli hafiza triage: Anthropic Haiku 4.5 (gunluk 03:30, max 20 batch)
30 3 * * * /opt/linux-ai-server/venv/bin/python3 /opt/linux-ai-server/scripts/memory-triage-llm.py
# Haftalık self-pentest — sahip olduğumuz tüm domain'leri tara (Pazar 03:30)
30 3 * * 0 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh self-pentest /opt/linux-ai-server/automation/self-pentest.sh
0 4 * * 0 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh nuclei-scan /opt/linux-ai-server/automation/nuclei-scan.sh
# DB retention — server.db + ci_tests.db (haftalik Pazar 02:00)
0 2 * * 0 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh db-retention /opt/linux-ai-server/automation/db-retention.sh
# Claude Code global npm guncelleme (haftalik Pazar 05:00)
0 5 * * 0 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh claude-code-update /opt/linux-ai-server/automation/claude-code-update.sh
# Otonom mod gunluk aktivite ozeti (her gun 06:00) — memory entry yazar, kullanici SessionStart'ta gorur
0 6 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh autonomous-summary /opt/linux-ai-server/automation/autonomous-daily-summary.sh
# Otonom mod health check (her 4 saat) — Ollama + Claude + DB + poller self-test
0 */4 * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh autonomous-health /opt/linux-ai-server/automation/autonomous-health-check.sh
# Otonom mod DLQ retry (her 15 dakika) — spawn_failures pending_retry'lari yeniden dene; 3x fail -> poison + Telegram
*/15 * * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh autonomous-retry /opt/linux-ai-server/automation/autonomous-spawn-retry.sh
# notify-cron: events tablosu pending (warn/critical) -> n8n bildirimi [FAZ3.2 step-d]
# DISABLED: Once (1) data/notify-cron.sh'i automation/'a tasiyip chmod +x yap,
#           (2) .env'e NOTIFY_CRON_ENABLED=true ekle, (3) bu satiri uncomment et + install-cron.sh calistir.
*/20 * * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh notify-cron /opt/linux-ai-server/automation/notify-cron.sh
# META-MONITOR: canli-sistem bekcilerini izle (liveness.check_all). Dead component ->
# DIRECT Telegram (dead-man's switch) + spine. */10. Bekcileri-kim-bekliyor katmani.
*/10 * * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh liveness-check /opt/linux-ai-server/automation/liveness-check.sh
# Renderhane bakiye uyarisi (panola-social VPS /api/health proxy, saatlik) — v2-06 PSOC-20260528
0 * * * * /opt/linux-ai-server/scripts/klipper-cron-wrap.sh renderhane-balance /opt/linux-ai-server/automation/social-renderhane-balance-alert.sh
# Haftalik VERI-ANALISTI (Pazartesi 08:00) — server.db + coverage.db trend, salt-okunur
# /claude (Max-plan) + db-query.sh -> discovery + Telegram ozet. Multi-uzman vizyon: 1/4.
# ENABLED 2026-06-06 (spot-check x3 PASS: rapor+discovery id537+Telegram). Opt-in kapi
# .env DATA_ANALYST_ENABLED=true; kapatmak icin flag=false (satir kalabilir, script no-op).
0 8 * * 1 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh data-analyst /opt/linux-ai-server/automation/data-analyst.sh
# Haftalik SEO denetimi (Pazartesi) — multi-uzman vizyon 3/4. Deterministik salt-okunur;
# bulgular ortak-hafizaya (type=bug -> SessionStart), Telegram yok.
# seo-gsc haftalik calismasi AYNI ZAMANDA OAuth refresh_token'i sicak tutar (6-ay-revoke onler)
# + token olurse OUTCOME:fail -> alert => ayri liveness-check gereksiz. ENABLED 2026-06-07.
15 8 * * 1 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh seo-audit /opt/linux-ai-server/automation/seo-audit.sh
30 8 * * 1 /opt/linux-ai-server/scripts/klipper-cron-wrap.sh seo-gsc /opt/linux-ai-server/automation/seo-gsc.sh
